# muninndb.service
[Unit]
Description=MuninnDB Memory Server
After=network.target

[Service]
Type=forking
# Data directory — override default (~/.muninn/data)
Environment=MUNINNDB_DATA=/opt/muninndb/data
# Network
Environment=MUNINN_LISTEN_HOST=0.0.0.0
Environment=MUNINN_CORS_ORIGINS=http://YOUR_IP:8476
# Auth
Environment=MUNINN_MCP_TOKEN=your-mcp-token-here
# Embeddings (choose one)
Environment=MUNINN_OPENAI_KEY=sk-...
# LLM Enrichment
Environment=MUNINN_ENRICH_URL=anthropic://claude-haiku-4-5-20251001
Environment=MUNINN_ANTHROPIC_KEY=sk-ant-...
ExecStart=/opt/muninndb-src/muninndb-server start
ExecStop=/opt/muninndb-src/muninndb-server stop
Restart=on-failure

[Install]
WantedBy=multi-user.target
